# traits/obj/persist/startup/windows.yaml
# Windows Startup Persistence Detection
# Detects Windows persistence mechanisms
# ATT&CK: T1547.001 (Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder)

defaults:
  attack: "T1547.001"
  crit: suspicious

traits:
  # === Startup folder atomic indicators ===
  - id: start-menu-startup
    desc: Start Menu Programs Startup path
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "\\Start Menu\\Programs\\Startup"

  - id: appdata-startup
    desc: AppData Roaming Startup path
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup"

  - id: shell-startup
    desc: shell:startup special folder
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "shell:startup"

  # === Registry Run key atomic indicators ===
  - id: software-run-key
    desc: SOFTWARE CurrentVersion Run path
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run"

  - id: hkcu-run
    desc: HKCU Run key reference
    crit: inert
    conf: 0.6
    platforms: [windows]
    if:
      type: string
      regex: "HKCU.{0,50}Run"

  - id: hklm-run
    desc: HKLM Run key reference
    crit: inert
    conf: 0.6
    platforms: [windows]
    if:
      type: string
      regex: "HKLM.{0,50}Run"

  - id: currentversion-run
    desc: CurrentVersion Run path
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "CurrentVersion\\Run"

  # === RunOnce atomic indicators ===
  - id: currentversion-runonce
    desc: CurrentVersion RunOnce path
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "CurrentVersion\\RunOnce"

  - id: runonce-word
    desc: RunOnce keyword
    crit: inert
    conf: 0.5
    platforms: [windows]
    if:
      type: string
      exact: "RunOnce"

  # === Python winreg atomic indicators ===
  - id: import-winreg
    desc: import winreg statement
    crit: inert
    conf: 0.5
    platforms: [windows]
    r#for: [python]
    if:
      type: string
      exact: "import winreg"

  - id: from-winreg
    desc: from winreg import
    crit: inert
    conf: 0.5
    platforms: [windows]
    r#for: [python]
    if:
      type: string
      exact: "from winreg"

  - id: underscore-winreg
    desc: _winreg module (Python 2)
    crit: inert
    conf: 0.5
    platforms: [windows]
    r#for: [python]
    if:
      type: string
      exact: "_winreg"

  - id: openkey-hkey
    desc: OpenKey HKEY call
    crit: inert
    conf: 0.6
    platforms: [windows]
    r#for: [python]
    if:
      type: string
      regex: "OpenKey.{0,30}HKEY"

  - id: setvalueex
    desc: SetValueEx call
    crit: inert
    conf: 0.6
    platforms: [windows]
    r#for: [python]
    if:
      type: string
      exact: "SetValueEx"

  - id: createkey
    desc: CreateKey call
    crit: inert
    conf: 0.6
    platforms: [windows]
    r#for: [python]
    if:
      type: string
      exact: "CreateKey"

  # === Scheduled tasks atomic indicators ===
  - id: schtasks-create-cmd
    desc: schtasks /create command
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      regex: "schtasks\\s*/create"

  - id: schtasks-exe-create
    desc: schtasks.exe create
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      regex: "schtasks\\.exe.{0,30}create"

  - id: itaskscheduler-iface
    desc: ITaskScheduler interface
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "ITaskScheduler"

  - id: taskscheduler-word
    desc: TaskScheduler keyword
    crit: inert
    conf: 0.6
    platforms: [windows]
    if:
      type: string
      exact: "TaskScheduler"

  # === Service installation atomic indicators ===
  - id: sc-create-cmd
    desc: sc create command
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      regex: "sc\\s+create"

  - id: createservice-api
    desc: CreateService API
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "CreateService"

  - id: installservice-fn
    desc: InstallService function
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "InstallService"

  - id: win32serviceutil-mod
    desc: win32serviceutil module
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "win32serviceutil"

  # === PowerShell registry atomic indicators ===
  - id: set-itemproperty-run
    desc: Set-ItemProperty Run
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      regex: "Set-ItemProperty.{0,50}Run"

  - id: new-itemproperty-run
    desc: New-ItemProperty Run
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      regex: "New-ItemProperty.{0,50}Run"

  - id: reg-add-run
    desc: reg add Run
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      regex: "reg\\s+add.{0,50}Run"

  # === WMI persistence atomic indicators ===
  - id: wmi-eventfilter
    desc: __EventFilter class
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "__EventFilter"

  - id: wmi-filtertoconsumer
    desc: __FilterToConsumerBinding class
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "__FilterToConsumerBinding"

  - id: wmi-commandlineconsumer
    desc: CommandLineEventConsumer class
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "CommandLineEventConsumer"

  - id: wmi-activescriptconsumer
    desc: ActiveScriptEventConsumer class
    crit: inert
    conf: 0.7
    platforms: [windows]
    if:
      type: string
      exact: "ActiveScriptEventConsumer"

  # === Executable payload atomic indicators ===
  - id: exe-extension
    desc: .exe file extension
    crit: inert
    conf: 0.4
    platforms: [windows]
    if:
      type: string
      exact: ".exe"

  - id: bat-extension
    desc: .bat file extension
    crit: inert
    conf: 0.4
    platforms: [windows]
    if:
      type: string
      exact: ".bat"

  - id: ps1-extension
    desc: .ps1 file extension
    crit: inert
    conf: 0.5
    platforms: [windows]
    if:
      type: string
      exact: ".ps1"

  - id: vbs-extension
    desc: .vbs file extension
    crit: inert
    conf: 0.5
    platforms: [windows]
    if:
      type: string
      exact: ".vbs"

  - id: powershell-word
    desc: powershell keyword
    crit: inert
    conf: 0.4
    platforms: [windows]
    if:
      type: string
      exact: "powershell"

  - id: cmd-exe
    desc: cmd.exe reference
    crit: inert
    conf: 0.4
    platforms: [windows]
    if:
      type: string
      exact: "cmd.exe"

composite_rules:
  # Startup folder composite
  - id: windows-startup-folder
    desc: "Targets Windows Startup folder for"
    conf: 0.85
    crit: suspicious
    platforms: [windows]
    count_min: 1
    any:
      - id: start-menu-startup
      - id: appdata-startup
      - id: shell-startup

  # Registry Run key composite
  - id: windows-registry-run
    desc: "Accesses Windows registry Run key"
    conf: 0.9
    crit: suspicious
    platforms: [windows]
    count_min: 1
    any:
      - id: software-run-key
      - id: hkcu-run
      - id: hklm-run
      - id: currentversion-run

  # RunOnce composite
  - id: windows-registry-runonce
    desc: "Accesses Windows registry RunOnce key"
    conf: 0.9
    crit: suspicious
    platforms: [windows]
    count_min: 1
    any:
      - id: currentversion-runonce
      - id: runonce-word

  # Python winreg import composite
  - id: windows-python-winreg-import
    desc: "Imports Python winreg module for"
    conf: 0.7
    crit: notable
    platforms: [windows]
    r#for: [python]
    count_min: 1
    any:
      - id: import-winreg
      - id: from-winreg
      - id: underscore-winreg

  # Python winreg composite
  - id: windows-python-winreg
    desc: "Accesses Windows registry via Python"
    conf: 0.8
    platforms: [windows]
    r#for: [python]
    count_min: 1
    any:
      - id: import-winreg
      - id: from-winreg
      - id: underscore-winreg
      - id: openkey-hkey
      - id: setvalueex
      - id: createkey

  # Registry Run access composite
  - id: windows-registry-run-access
    desc: "Accesses Windows registry Run key"
    conf: 0.8
    crit: notable
    platforms: [windows]
    count_min: 1
    any:
      - id: currentversion-run
      - id: openkey-hkey

  # Scheduled tasks composite
  - id: windows-schtasks
    desc: "Creates scheduled task for persistence"
    conf: 0.85
    crit: suspicious
    attack: "T1053.005"
    platforms: [windows]
    count_min: 1
    any:
      - id: schtasks-create-cmd
      - id: schtasks-exe-create
      - id: itaskscheduler-iface
      - id: taskscheduler-word

  # Scheduled task create composite
  - id: windows-schtasks-create
    desc: "Creates Windows scheduled task"
    conf: 0.8
    crit: notable
    attack: "T1053.005"
    platforms: [windows]
    count_min: 1
    any:
      - id: schtasks-create-cmd
      - id: schtasks-exe-create
      - id: itaskscheduler-iface

  # Service installation composite
  - id: windows-service-install
    desc: "Installs Windows service for persistence"
    conf: 0.85
    attack: "T1543.003"
    platforms: [windows]
    count_min: 1
    any:
      - id: sc-create-cmd
      - id: createservice-api
      - id: installservice-fn
      - id: win32serviceutil-mod

  # PowerShell registry composite
  - id: windows-powershell-registry
    desc: "Modifies registry for persistence via"
    conf: 0.9
    crit: suspicious
    platforms: [windows]
    count_min: 1
    any:
      - id: set-itemproperty-run
      - id: new-itemproperty-run
      - id: reg-add-run

  # WMI persistence composite
  - id: windows-wmi-persistence
    desc: "Sets up WMI event subscription"
    conf: 0.9
    crit: suspicious
    attack: "T1546.003"
    platforms: [windows]
    count_min: 1
    any:
      - id: wmi-eventfilter
      - id: wmi-filtertoconsumer
      - id: wmi-commandlineconsumer
      - id: wmi-activescriptconsumer

  # Executable payload composite
  - id: windows-executable-payload
    desc: "References executable file types or"
    conf: 0.6
    crit: notable
    platforms: [windows]
    count_min: 1
    any:
      - id: exe-extension
      - id: bat-extension
      - id: ps1-extension
      - id: vbs-extension
      - id: powershell-word
      - id: cmd-exe

  # Python registry persistence pattern
  - id: windows-python-reg-persistence
    desc: "Python Windows registry persistence"
    crit: suspicious
    conf: 0.95
    attack: "T1547.001"
    mbc: "F0006"
    platforms: [windows]
    r#for: [python]
    all:
      - id: windows-python-winreg-import
      - id: windows-registry-run-access

  # Scheduled task with suspicious payload
  - id: windows-schtasks-payload
    desc: "Scheduled task with executable payload"
    crit: suspicious
    conf: 0.95
    attack: "T1053.005"
    mbc: "F0006"
    platforms: [windows]
    all:
      - id: windows-schtasks-create
      - id: windows-executable-payload
