# traits/persist/callback/polling/javascript.yaml
# Persistent Polling - JavaScript/Node.js
# Detects setInterval/setTimeout patterns used for persistent callbacks
# ATT&CK: T1053 (Scheduled Task/Job)

traits:
  # === setInterval HTTP atomic indicators ===
  - id: setInterval-http-call
    description: setInterval with http
    crit: inert
    confidence: 0.6
    file_types: [javascript]
    condition:
      type: string
      regex: "setInterval.{0,100}http"

  - id: setInterval-axios
    description: setInterval with axios
    crit: inert
    confidence: 0.6
    file_types: [javascript]
    condition:
      type: string
      regex: "setInterval.{0,100}axios"

  - id: setInterval-fetch
    description: setInterval with fetch
    crit: inert
    confidence: 0.6
    file_types: [javascript]
    condition:
      type: string
      regex: "setInterval.{0,100}fetch"

  - id: setInterval-function
    description: setInterval with function callback
    crit: inert
    confidence: 0.7
    file_types: [javascript]
    condition:
      type: string
      regex: "setInterval\\s*\\(\\s*function"

  - id: setTimeout-recursive
    description: setTimeout in recursive pattern
    crit: notable
    confidence: 0.7
    file_types: [javascript]
    condition:
      type: string
      regex: "setTimeout.{0,100}setTimeout"

composite_rules:
  - id: setInterval-http
    description: setInterval with HTTP calls (persistent polling)
    crit: suspicious
    confidence: 0.85
    attack: "T1053"
    file_types: [javascript]
    count: 1
    any:
      - id: setInterval-http-call
      - id: setInterval-axios
      - id: setInterval-fetch
