# traits/persist/callback/polling/javascript.yaml
# Persistent Polling - JavaScript/Node.js
# Detects setInterval/setTimeout patterns used for persistent callbacks
# ATT&CK: T1053 (Scheduled Task/Job)

traits:
  # === setInterval HTTP atomic indicators ===
  - id: setInterval-http-call
    description: setInterval with http
    criticality: inert
    confidence: 0.6
    file_types: [javascript]
    condition:
      type: string
      regex: "setInterval.*http"

  - id: setInterval-axios
    description: setInterval with axios
    criticality: inert
    confidence: 0.6
    file_types: [javascript]
    condition:
      type: string
      regex: "setInterval.*axios"

  - id: setInterval-fetch
    description: setInterval with fetch
    criticality: inert
    confidence: 0.6
    file_types: [javascript]
    condition:
      type: string
      regex: "setInterval.*fetch"

  - id: setInterval-function
    description: setInterval with function callback
    criticality: inert
    confidence: 0.7
    file_types: [javascript]
    condition:
      type: string
      regex: "setInterval\\s*\\(\\s*function"

  - id: setTimeout-recursive
    description: setTimeout in recursive pattern
    criticality: notable
    confidence: 0.7
    file_types: [javascript]
    condition:
      type: string
      regex: "setTimeout.*setTimeout"

composite_rules:
  - id: setInterval-http
    description: setInterval with HTTP calls (persistent polling)
    criticality: suspicious
    confidence: 0.85
    attack: "T1053"
    file_types: [javascript]
    requires_count: 1
    conditions:
      - id: setInterval-http-call
      - id: setInterval-axios
      - id: setInterval-fetch
