# traits/persistence/startup/windows.yaml
# Windows Startup Persistence Detection
# Detects Windows persistence mechanisms
# ATT&CK: T1547.001 (Boot or Logon Autostart Execution: Registry Run Keys / Startup Folder)

defaults:
  attack: "T1547.001"
  criticality: suspicious

traits:
  # Windows Startup folder paths
  - id: persistence/startup/windows-startup-folder
    description: "Windows Startup folder path"
    confidence: 0.85
    criticality: hostile
    platforms: [windows]
    condition:
      type: string
      regex: "\\\\Start Menu\\\\Programs\\\\Startup|AppData\\\\Roaming\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup|shell:startup"

  # Registry Run keys (common persistence location)
  - id: persistence/startup/windows-registry-run
    description: "Windows registry Run key access"
    confidence: 0.9
    criticality: hostile
    platforms: [windows]
    condition:
      type: string
      regex: "SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run|HKCU.*Run|HKLM.*Run|CurrentVersion\\\\Run"

  # Registry RunOnce keys
  - id: persistence/startup/windows-registry-runonce
    description: "Windows registry RunOnce key access"
    confidence: 0.9
    criticality: hostile
    platforms: [windows]
    condition:
      type: string
      regex: "CurrentVersion\\\\RunOnce|RunOnce"

  # winreg module (Python Windows registry access)
  - id: persistence/startup/windows-python-winreg
    description: "Windows registry access via Python winreg"
    confidence: 0.8
    platforms: [windows]
    file_types: [python]
    condition:
      type: string
      regex: "import winreg|from winreg|_winreg|OpenKey.*HKEY|SetValueEx|CreateKey"

  # Scheduled tasks (schtasks)
  - id: persistence/startup/windows-schtasks
    description: "Windows scheduled task creation"
    confidence: 0.85
    criticality: hostile
    attack: "T1053.005"
    platforms: [windows]
    condition:
      type: string
      regex: "schtasks\\s*/create|schtasks\\.exe.*create|ITaskScheduler|TaskScheduler"

  # Windows service installation
  - id: persistence/startup/windows-service-install
    description: "Windows service installation"
    confidence: 0.85
    attack: "T1543.003"
    platforms: [windows]
    condition:
      type: string
      regex: "sc\\s+create|CreateService|InstallService|win32serviceutil"

  # PowerShell registry persistence
  - id: persistence/startup/windows-powershell-registry
    description: "PowerShell registry modification for persistence"
    confidence: 0.9
    criticality: hostile
    platforms: [windows]
    condition:
      type: string
      regex: "Set-ItemProperty.*Run|New-ItemProperty.*Run|reg\\s+add.*Run"

  # WMI event subscription persistence
  - id: persistence/startup/windows-wmi-persistence
    description: "WMI event subscription persistence"
    confidence: 0.9
    criticality: hostile
    attack: "T1546.003"
    platforms: [windows]
    condition:
      type: string
      regex: "__EventFilter|__FilterToConsumerBinding|CommandLineEventConsumer|ActiveScriptEventConsumer"

composite_rules:
  # Python registry persistence pattern
  - id: persistence/startup/windows-python-reg-persistence
    description: "Python Windows registry persistence"
    criticality: hostile
    confidence: 0.95
    platforms: [windows]
    file_types: [python]
    requires_all:
      - type: string
        regex: "import winreg|from winreg|_winreg"
      - type: string
        regex: "CurrentVersion\\\\Run|OpenKey.*HKEY.*Run"

  # Scheduled task with suspicious payload
  - id: persistence/startup/windows-schtasks-payload
    description: "Scheduled task with executable payload"
    criticality: hostile
    confidence: 0.95
    attack: "T1053.005"
    platforms: [windows]
    requires_all:
      - type: string
        regex: "schtasks.*create|ITaskScheduler"
      - type: string
        regex: "\\.exe|\\.bat|\\.ps1|\\.vbs|powershell|cmd\\.exe"
