# traits/known-tools/backdoor/upstyle/traits.yaml
# PAN-OS Upstyle Backdoor Detection
# ATT&CK: T1505 (Server Software Component)
# Reference: Palo Alto Networks firewall backdoor

defaults:
  file_types: [python, elf, macho, pe, so, dylib, dll, shell, javascript]

traits:
  # === Backdoor Name Variants ===
  - id: upstyle-lower
    desc: upstyle malware name (lowercase)
    crit: suspicious
    conf: 0.95
    if:
      type: string
      exact: "upstyle"

  - id: upstyle-upper
    desc: UPSTYLE malware name (uppercase)
    crit: suspicious
    conf: 0.95
    if:
      type: string
      exact: "UPSTYLE"

  - id: upstyle-title
    desc: Upstyle malware name (title case)
    crit: suspicious
    conf: 0.95
    if:
      type: string
      exact: "Upstyle"

  # === PAN-OS Paths ===
  - id: var-appweb-path
    desc: PAN-OS /var/appweb path
    crit: inert
    conf: 0.75
    if:
      type: string
      exact: "/var/appweb"

  - id: opt-panlogs-path
    desc: PAN-OS /opt/panlogs path
    crit: inert
    conf: 0.75
    if:
      type: string
      exact: "/opt/panlogs"

  - id: globalprotect-string
    desc: GlobalProtect VPN reference
    crit: inert
    conf: 0.7
    if:
      type: string
      exact: "GlobalProtect"

  - id: opt-pancfg-path
    desc: PAN-OS /opt/pancfg path
    crit: inert
    conf: 0.75
    if:
      type: string
      exact: "/opt/pancfg"

  # === Backdoor Markers ===
  - id: update-py-file
    desc: update.py backdoor filename
    crit: inert
    conf: 0.75
    if:
      type: string
      exact: "update.py"

  - id: system-upstyle-string
    desc: system_upstyle backdoor marker
    crit: suspicious
    conf: 0.95
    if:
      type: string
      exact: "system_upstyle"

  # === Python Execution ===
  - id: subprocess-import
    desc: subprocess Python module
    crit: inert
    conf: 0.6
    if:
      type: string
      exact: "subprocess"

  - id: os-system-call
    desc: os.system Python call
    crit: inert
    conf: 0.65
    if:
      type: string
      exact: "os.system"

  - id: exec-call
    desc: exec() Python call
    crit: inert
    conf: 0.7
    if:
      type: string
      exact: "exec("

  # === PAN-OS Context (case variants) ===
  - id: panos-lower
    desc: panos string (lowercase)
    crit: inert
    conf: 0.7
    if:
      type: string
      exact: "panos"

  - id: panos-upper
    desc: PANOS string (uppercase)
    crit: inert
    conf: 0.7
    if:
      type: string
      exact: "PANOS"

  - id: panos-title
    desc: PanOS string (title case)
    crit: inert
    conf: 0.7
    if:
      type: string
      exact: "PanOS"

  # === Web Shell Behavior ===
  - id: request-get
    desc: request.GET web parameter
    crit: inert
    conf: 0.65
    if:
      type: string
      exact: "request.GET"

  - id: request-post
    desc: request.POST web parameter
    crit: inert
    conf: 0.65
    if:
      type: string
      exact: "request.POST"

  # === Persistence Mechanisms ===
  - id: etc-cron-path
    desc: /etc/cron directory path
    crit: inert
    conf: 0.6
    if:
      type: string
      exact: "/etc/cron"

  - id: crontab-cmd
    desc: crontab command
    crit: inert
    conf: 0.6
    if:
      type: string
      exact: "crontab"

  - id: systemctl-cmd
    desc: systemctl command
    crit: inert
    conf: 0.6
    if:
      type: string
      exact: "systemctl"

  - id: etc-init-d-path
    desc: /etc/init.d directory path
    crit: inert
    conf: 0.6
    if:
      type: string
      exact: "/etc/init.d"

  # === PAN-OS Services ===
  - id: mgmtsrvr-service
    desc: PAN-OS mgmtsrvr service
    crit: inert
    conf: 0.75
    if:
      type: string
      exact: "mgmtsrvr"

  - id: appweb-service
    desc: PAN-OS appweb service
    crit: inert
    conf: 0.75
    if:
      type: string
      exact: "appweb"

composite_rules:
  # === Helper composites ===
  - id: upstyle-names
    description: Upstyle name (case variants)
    any:
      - id: upstyle-lower
      - id: upstyle-upper
      - id: upstyle-title

  - id: pan-os-paths
    description: PAN-OS specific paths
    any:
      - id: var-appweb-path
      - id: opt-panlogs-path
      - id: globalprotect-string
      - id: opt-pancfg-path

  - id: backdoor-markers
    description: Backdoor file markers
    any:
      - id: update-py-file
      - id: system-upstyle-string

  - id: python-execution
    description: Python execution methods
    any:
      - id: subprocess-import
      - id: os-system-call
      - id: exec-call

  - id: panos-context
    description: PAN-OS context strings
    any:
      - id: panos-lower
      - id: panos-upper
      - id: panos-title

  - id: web-params
    description: Web request parameters
    any:
      - id: request-get
      - id: request-post

  - id: cron-persistence
    description: Cron-based persistence
    any:
      - id: etc-cron-path
      - id: crontab-cmd

  - id: service-persistence
    description: Service-based persistence
    any:
      - id: systemctl-cmd
      - id: etc-init-d-path

  - id: pan-services
    description: PAN-OS services
    any:
      - id: mgmtsrvr-service
      - id: appweb-service

  - id: pan-with-backdoor-markers
    description: PAN-OS paths with backdoor markers
    all:
      - id: pan-os-paths
      - id: backdoor-markers

  # === Migrated from YARA ===
  - id: backdoor/upstyle/generic
    description: PAN-OS Upstyle backdoor patterns
    crit: suspicious
    conf: 0.9
    attack: "T1505"
    any:
      - id: upstyle-names
      - id: pan-with-backdoor-markers

  - id: panos-with-python
    description: PAN-OS context with Python execution
    all:
      - id: panos-context
      - id: python-execution

  - id: panos-with-web
    description: PAN-OS context with web parameters
    all:
      - id: panos-context
      - id: web-params

  - id: backdoor/upstyle/cmd-exec
    description: Upstyle command execution capability
    crit: suspicious
    conf: 0.85
    attack: "T1059"
    any:
      - id: upstyle-names
      - id: panos-with-python
      - id: panos-with-web

  - id: cron-or-service
    description: Cron or service persistence
    any:
      - id: cron-persistence
      - id: service-persistence

  - id: pan-services-persist
    description: PAN-OS services with persistence
    all:
      - id: pan-services
      - id: cron-or-service

  - id: backdoor/upstyle/persistence
    description: Upstyle persistence mechanism
    crit: suspicious
    conf: 0.85
    attack: "T1546"
    any:
      - id: upstyle-names
      - id: pan-services-persist

  - id: backdoor/upstyle/detected
    description: Upstyle backdoor detected
    crit: suspicious
    conf: 0.95
