# traits/known-tools/remote-access-trojan/pupy/traits.yaml
# Pupy RAT Detection
# ATT&CK: T1059 (Command and Scripting Interpreter)
# MBC: B0030 (Remote Access)

traits:
  # Pupy reflective DLL/SO loader
  - id: c2-framework/pupy/reflective-loader
    description: Pupy reflective loader patterns
    crit: suspicious
    conf: 0.9
    mbc: "B0030"
    attack: "T1620"
    file_types: [elf, macho, pe, so, dylib, dll]
    condition:
      type: yara
      source: |
        rule pupy_reflective {
          strings:
            $pupy1 = "pupy" ascii nocase
            $pupy2 = "pupyimporter" ascii
            $ref1 = "reflective" ascii nocase
            $ref2 = "memimporter" ascii
            $ref3 = "MemoryModule" ascii
            $py1 = "Py_Initialize" ascii
            $py2 = "PyRun_SimpleString" ascii
            $py3 = "PyEval_" ascii
          condition:
            (any of ($pupy*) and any of ($py*)) or (2 of ($ref*) and any of ($py*))
        }

  # Pupy Python stager
  - id: c2-framework/pupy/python-stager
    description: Pupy Python stager patterns
    crit: suspicious
    conf: 0.85
    mbc: "B0030"
    attack: "T1059.006"
    file_types: [python]
    condition:
      type: yara
      source: |
        rule pupy_python {
          strings:
            $import1 = "import pupy" ascii
            $import2 = "from pupy" ascii
            $import3 = "pupylib" ascii
            $loader1 = "PupyClient" ascii
            $loader2 = "pupymodules" ascii
            $loader3 = "pupygen" ascii
            $transport1 = "obfs3" ascii
            $transport2 = "scramblesuit" ascii
            $transport3 = "ssl_rsa" ascii
          condition:
            (any of ($import*) and any of ($loader*)) or (any of ($loader*) and any of ($transport*))
        }

  # Pupy rpyc transport
  - id: c2-framework/pupy/rpyc-transport
    description: Pupy rpyc-based transport layer
    crit: suspicious
    conf: 0.8
    mbc: "B0030"
    attack: "T1071"
    file_types: [elf, macho, pe, so, dylib, dll, python]
    condition:
      type: yara
      source: |
        rule pupy_rpyc {
          strings:
            $rpyc1 = "rpyc" ascii
            $rpyc2 = "RPyC" ascii
            $rpyc3 = "brine" ascii
            $stream1 = "PupyChannel" ascii
            $stream2 = "PupyStream" ascii
            $proto1 = "PupyService" ascii
            $proto2 = "PupyConnection" ascii
          condition:
            (any of ($rpyc*) and any of ($stream*)) or (any of ($stream*) and any of ($proto*))
        }

  # Pupy in-memory module loading
  - id: c2-framework/pupy/inmemory
    description: Pupy in-memory Python module loading
    crit: suspicious
    conf: 0.85
    mbc: "B0030"
    attack: "T1620"
    file_types: [elf, macho, pe, so, dylib, dll]
    condition:
      type: yara
      source: |
        rule pupy_inmemory {
          strings:
            $mem1 = "memimport" ascii
            $mem2 = "MemoryModule" ascii
            $mem3 = "memfd_create" ascii
            $zip1 = "zipimporter" ascii
            $py1 = "marshal" ascii
            $py2 = "importlib" ascii
            // Pupy-specific
            $pupy1 = "pupy" ascii nocase
            $pupy2 = "pupylib" ascii
          condition:
            (any of ($mem1, $mem2) and any of ($py*)) or
            (any of ($pupy*) and any of ($mem*, $zip*, $py*))
        }

  # Pupy persistence mechanisms
  - id: c2-framework/pupy/persistence
    description: Pupy persistence module patterns
    crit: suspicious
    conf: 0.8
    mbc: "B0030"
    attack: "T1547"
    file_types: [python, shell]
    condition:
      type: yara
      source: |
        rule pupy_persistence {
          strings:
            $persist1 = "persistence" ascii
            $persist2 = "autorun" ascii
            $cron1 = "/etc/cron" ascii
            $cron2 = "crontab" ascii
            $systemd1 = "/etc/systemd" ascii
            $systemd2 = ".service" ascii
            $rc1 = "/etc/rc.local" ascii
            $rc2 = "/etc/init.d" ascii
          condition:
            any of ($persist*) and (2 of ($cron*, $systemd*, $rc*))
        }

  # Pupy embedded Python
  - id: c2-framework/pupy/embedded-python
    description: Pupy embedded Python interpreter
    crit: notable
    conf: 0.75
    mbc: "B0030"
    attack: "T1059.006"
    file_types: [elf, macho, pe, so, dylib, dll]
    condition:
      type: yara
      source: |
        rule pupy_embedded_python {
          strings:
            $py_api1 = "Py_SetProgramName" ascii
            $py_api2 = "Py_SetPythonHome" ascii
            $py_api3 = "Py_InitializeEx" ascii
            $py_api4 = "PyGILState_" ascii
            $frozen1 = "FrozenImporter" ascii
            $frozen2 = "_frozen_importlib" ascii
            // High entropy section (encrypted payload)
            $magic = { 50 4B 03 04 }  // ZIP magic
          condition:
            (3 of ($py_api*) and any of ($frozen*)) or ($magic and any of ($py_api*))
        }

composite_rules:
  # Full Pupy RAT detection
  - id: c2-framework/pupy/detected
    description: Pupy RAT detected
    crit: suspicious
    conf: 0.95
    mbc: "B0030"
    attack: "T1071"
    file_types: [elf, macho, pe, so, dylib, dll, python]
    any:
      - id: c2-framework/pupy/reflective-loader
      - id: c2-framework/pupy/python-stager
      - id: c2-framework/pupy/rpyc-transport
      - id: c2-framework/pupy/inmemory
