# traits/persist/autostart/linux.yaml
# Linux Desktop Autostart Persistence
# Detects XDG autostart and other Linux GUI persistence mechanisms
# ATT&CK: T1547.013 (Boot or Logon Autostart Execution: XDG Autostart Entries)

defaults:
  attack: "T1547.013"
  crit: suspicious
  platforms: [linux]

traits:
  # === XDG autostart atomic indicators ===
  - id: config-autostart
    description: .config/autostart directory
    crit: inert
    confidence: 0.7
    condition:
      type: string
      exact: ".config/autostart"

  - id: etc-xdg-autostart
    description: /etc/xdg/autostart directory
    crit: inert
    confidence: 0.7
    condition:
      type: string
      exact: "/etc/xdg/autostart"

  - id: xdg-config-autostart
    description: XDG_CONFIG_HOME autostart
    crit: inert
    confidence: 0.6
    condition:
      type: string
      regex: "XDG_CONFIG_HOME.{0,50}autostart"

  # === Desktop entry atomic indicators ===
  - id: dot-desktop-ext
    description: .desktop file extension
    crit: inert
    confidence: 0.5
    condition:
      type: string
      exact: ".desktop"

  - id: desktop-entry-header
    description: "[Desktop Entry] header"
    crit: inert
    confidence: 0.6
    condition:
      type: string
      exact: "[Desktop Entry]"

  - id: x-gnome-autostart
    description: X-GNOME-Autostart directive
    crit: inert
    confidence: 0.7
    condition:
      type: string
      exact: "X-GNOME-Autostart"

  # === User systemd atomic indicators ===
  - id: config-systemd-user
    description: .config/systemd/user directory
    crit: inert
    confidence: 0.7
    condition:
      type: string
      exact: ".config/systemd/user"

  - id: local-share-systemd-user
    description: .local/share/systemd/user directory
    crit: inert
    confidence: 0.7
    condition:
      type: string
      exact: ".local/share/systemd/user"

  - id: systemctl-user-enable
    description: systemctl --user enable
    crit: inert
    confidence: 0.7
    condition:
      type: string
      regex: "systemctl.{0,30}--user.{0,30}enable"

  # === rc.local atomic indicators ===
  - id: etc-rc-local
    description: /etc/rc.local file
    crit: inert
    confidence: 0.8
    attack: "T1037.004"
    condition:
      type: string
      exact: "/etc/rc.local"

  - id: etc-rc-d-rc-local
    description: /etc/rc.d/rc.local file
    crit: inert
    confidence: 0.8
    attack: "T1037.004"
    condition:
      type: string
      exact: "/etc/rc.d/rc.local"

  # === init.d atomic indicators ===
  - id: etc-init-d
    description: /etc/init.d/ directory
    crit: inert
    confidence: 0.6
    attack: "T1037.004"
    file_types: [shell, elf, macho, python]
    condition:
      type: string
      exact: "/etc/init.d/"

  - id: update-rc-d-cmd
    description: update-rc.d command
    crit: inert
    confidence: 0.7
    attack: "T1037.004"
    file_types: [shell, elf, macho, python]
    condition:
      type: string
      exact: "update-rc.d"

  - id: chkconfig-cmd
    description: chkconfig command
    crit: inert
    confidence: 0.7
    attack: "T1037.004"
    file_types: [shell, elf, macho, python]
    condition:
      type: string
      exact: "chkconfig"

  # === Desktop entry exec atomic indicators ===
  - id: exec-eq-directive
    description: Exec= directive
    crit: inert
    confidence: 0.6
    condition:
      type: string
      exact: "Exec="

  - id: hidden-false-directive
    description: Hidden=false directive
    crit: inert
    confidence: 0.5
    condition:
      type: string
      exact: "Hidden=false"

  - id: type-application-directive
    description: Type=Application directive
    crit: inert
    confidence: 0.5
    condition:
      type: string
      exact: "Type=Application"

composite_rules:
  # XDG autostart composite
  - id: linux-xdg-autostart
    description: "XDG autostart directory access"
    confidence: 0.85
    crit: suspicious
    count_min: 1
    any:
      - id: config-autostart
      - id: etc-xdg-autostart
      - id: xdg-config-autostart

  # Desktop entry composite
  - id: linux-desktop-entry
    description: "Desktop entry file creation (.desktop)"
    confidence: 0.8
    count_min: 1
    any:
      - id: dot-desktop-ext
      - id: desktop-entry-header
      - id: x-gnome-autostart

  # User systemd composite
  - id: linux-user-systemd
    description: "User-level systemd service persistence"
    confidence: 0.85
    attack: "T1543.002"
    count_min: 1
    any:
      - id: config-systemd-user
      - id: local-share-systemd-user
      - id: systemctl-user-enable

  # rc.local composite
  - id: linux-rc-local
    description: "rc.local persistence mechanism"
    confidence: 0.9
    crit: suspicious
    attack: "T1037.004"
    count_min: 1
    any:
      - id: etc-rc-local
      - id: etc-rc-d-rc-local

  # init.d composite
  - id: linux-init-script
    description: "init.d script persistence"
    confidence: 0.85
    attack: "T1037.004"
    file_types: [shell, elf, macho, python]
    count_min: 1
    any:
      - id: etc-init-d
      - id: update-rc-d-cmd
      - id: chkconfig-cmd

  # XDG autostart directory composite
  - id: linux-xdg-autostart-dir
    description: "XDG autostart directory path"
    confidence: 0.8
    crit: notable
    count_min: 1
    any:
      - id: config-autostart
      - id: etc-xdg-autostart

  # Desktop entry exec composite
  - id: linux-desktop-entry-exec
    description: "Desktop entry executable directive"
    confidence: 0.8
    crit: notable
    count_min: 1
    any:
      - id: exec-eq-directive
      - id: hidden-false-directive
      - id: type-application-directive

  # XDG autostart with executable
  - id: linux-xdg-executable
    description: "XDG autostart with executable command"
    crit: suspicious
    confidence: 0.95
    attack: "T1547.013"
    mbc: "F0006"
    file_types: [shell, python, elf]
    all:
      - id: linux-xdg-autostart-dir
      - id: linux-desktop-entry-exec
