# traits/persist/autostart/linux.yaml
# Linux Desktop Autostart Persistence
# Detects XDG autostart and other Linux GUI persistence mechanisms
# ATT&CK: T1547.013 (Boot or Logon Autostart Execution: XDG Autostart Entries)

defaults:
  attack: "T1547.013"
  criticality: suspicious
  platforms: [linux]

traits:
  # XDG autostart directories
  - id: linux-xdg-autostart
    description: "XDG autostart directory access"
    confidence: 0.85
    criticality: suspicious
    condition:
      type: string
      regex: "\\.config/autostart|/etc/xdg/autostart|XDG_CONFIG_HOME.*autostart"

  # Desktop entry files
  - id: linux-desktop-entry
    description: "Desktop entry file creation (.desktop)"
    confidence: 0.8
    condition:
      type: string
      regex: "\\.desktop|\\[Desktop Entry\\]|X-GNOME-Autostart"

  # User systemd service persistence
  - id: linux-user-systemd
    description: "User-level systemd service persistence"
    confidence: 0.85
    attack: "T1543.002"
    condition:
      type: string
      regex: "\\.config/systemd/user|\\.local/share/systemd/user|systemctl.*--user.*enable"

  # rc.local persistence
  - id: linux-rc-local
    description: "rc.local persistence mechanism"
    confidence: 0.9
    criticality: suspicious
    attack: "T1037.004"
    condition:
      type: string
      regex: "/etc/rc\\.local|/etc/rc\\.d/rc\\.local"

  # init.d script persistence
  - id: linux-init-script
    description: "init.d script persistence"
    confidence: 0.85
    attack: "T1037.004"
    # Exclude JS - syntax highlighters have shell command keywords
    file_types: [shell, elf, macho, python]
    condition:
      type: string
      regex: "/etc/init\\.d/|update-rc\\.d|chkconfig"

composite_rules:
  # XDG autostart with executable
  - id: linux-xdg-executable
    description: "XDG autostart with executable command"
    criticality: suspicious
    confidence: 0.95
    attack: "T1547.013"
    mbc: "F0006"
    file_types: [shell, python, elf]
    requires_all:
      - type: string
        regex: "\\.config/autostart|/etc/xdg/autostart"
      - type: string
        regex: "Exec=|Hidden=false|Type=Application"
