# traits/obj/persist/botnet/mirai/traits.yaml
# Mirai Complete Persistence Chain
# Detects complete Mirai persistence strategy combining multiple techniques
# ATT&CK: T1547 (Boot or Logon Autostart Execution)
# MBC: B0014 (Persistence)

traits:
  # === Individual Persistence Components ===
  - id: mirai-persistent-service
    desc: "Persistent service"
    crit: suspicious
    conf: 0.8
    for: [elf, shell]
    if:
      type: string
      any:
        - substr: "/usr/lib/systemd/systemd"
        - substr: "systemd"
    attack: "T1547"

  - id: mirai-backdoor-daemons
    desc: "Backdoor daemon processes"
    crit: suspicious
    conf: 0.85
    for: [elf, shell]
    if:
      type: string
      any:
        - substr: "telnetd"
        - substr: "sshd"
        - substr: "dropbear"
        - substr: "httpd"
    attack: "T1021.004"

  - id: mirai-process-monitoring
    desc: "Process monitoring"
    crit: suspicious
    conf: 0.85
    for: [elf, shell]
    if:
      type: string
      substr: "/dev/watchdog"
    attack: "T1547"

  - id: mirai-malware-competition
    desc: "Competing malware elimination"
    crit: suspicious
    conf: 0.9
    for: [elf, shell]
    if:
      type: string
      any:
        - substr: "SHINJI"
        - substr: "/Killer"
    attack: "T1531"

  # === Persistence Activation ===
  - id: mirai-enable-services
    desc: "Enable services on startup"
    crit: suspicious
    conf: 0.8
    for: [elf, shell]
    if:
      type: string
      any:
        - regex: "systemctl\\s+enable"
        - regex: "chkconfig\\s+.*on"
        - regex: "update-rc\\.d\\s+.*enable"
    attack: "T1547"

composite_rules:
  # Layer 1 handled directly by mirai-backdoor-daemons trait

  - id: mirai-persistence-layer-2
    desc: "Mirai service persistence"
    crit: suspicious
    conf: 0.9
    for: [elf, shell]
    attack: "T1547"
    any:
      - id: mirai-persistent-service
      - id: mirai-enable-services

  - id: mirai-persistence-layer-3
    desc: "Mirai monitoring and competition"
    crit: suspicious
    conf: 0.9
    for: [elf, shell]
    attack: "T1547"
    all:
      - id: mirai-process-monitoring
      - id: mirai-malware-competition

  - id: mirai-complete-persistence-chain
    desc: "Mirai complete persistence"
    crit: hostile
    conf: 0.95
    for: [elf, shell]
    attack: "T1547"
    mbc: "B0014"
    all:
      - id: mirai-backdoor-daemons
      - id: mirai-persistence-layer-2
      - id: mirai-persistence-layer-3

  - id: mirai-sophisticated-resilience
    desc: "Mirai sophisticated resilience"
    crit: hostile
    conf: 0.95
    for: [elf, shell]
    attack: "T1547"
    mbc: "B0014"
    all:
      - id: known/malware/botnet/mirai
      - id: mirai-complete-persistence-chain
