# traits/obj/lateral/supply-chain/vscode/malicious.yaml
# Malicious VSCode Extension Patterns
# Detects high-confidence malicious patterns in VSCode extensions
# These combine multiple indicators to reduce false positives

traits:
  # NOTE: YARA kept for api-usage - requires fullword match on "activate" which
  # DISSECT doesn't support, plus regex patterns for require/import
  - id: api-usage
    desc: Uses VSCode API (extension marker)
    crit: inert
    conf: 1.0
    r#for: [javascript, typescript]
    if:
      type: yara
      source: |
        rule vscode_api_usage {
          strings:
            $require = /require\(['"]vscode['"]\)/ ascii
            $import = /from\s+['"]vscode['"]/ ascii
            $namespace = "vscode." ascii
            $activate = "activate" ascii fullword
          condition:
            any of them
        }

  # === Atomic indicators for post-install-exec ===
  - id: activate-word
    desc: activate keyword
    crit: inert
    conf: 0.7
    r#for: [javascript, typescript]
    if:
      type: string
      substr: "activate"

  - id: globalState-word
    desc: globalState keyword
    crit: inert
    conf: 0.7
    r#for: [javascript, typescript]
    if:
      type: string
      substr: "globalState"

  - id: postInstall-word
    desc: postInstall keyword
    crit: inert
    conf: 0.7
    r#for: [javascript, typescript]
    if:
      type: string
      substr: "postInstall"

  - id: exec-pattern
    desc: exec function pattern
    crit: inert
    conf: 0.7
    r#for: [javascript, typescript]
    if:
      type: string
      regex: "exec(Sync|Async)?"

  # === Atomic indicators for credential-theft ===
  - id: ssh-dir
    desc: .ssh directory reference
    crit: notable
    conf: 0.7
    r#for: [javascript, typescript]
    if:
      type: string
      exact: ".ssh"

  - id: aws-dir
    desc: .aws directory reference
    crit: notable
    conf: 0.7
    r#for: [javascript, typescript]
    if:
      type: string
      exact: ".aws"

  - id: npmrc-file
    desc: .npmrc file reference
    crit: notable
    conf: 0.7
    r#for: [javascript, typescript]
    if:
      type: string
      substr: ".npmrc"

  - id: env-file
    desc: .env file reference
    crit: notable
    conf: 0.6
    r#for: [javascript, typescript]
    if:
      type: string
      exact: ".env"

  - id: readFile-pattern
    desc: readFile function pattern
    crit: inert
    conf: 0.6
    r#for: [javascript, typescript]
    if:
      type: string
      regex: "read(File|FileSync)"

  - id: post-method
    desc: POST HTTP method
    crit: inert
    conf: 0.6
    r#for: [javascript, typescript]
    if:
      type: string
      exact: "POST"

  # VSCode namespace usage pattern
  - id: vscode-namespace
    desc: Uses VSCode namespace
    crit: inert
    conf: 1.0
    r#for: [javascript, typescript]
    if:
      type: string
      regex: "vscode\\."

  # Persistence location references
  - id: persistence-locations
    desc: References persistence locations (startup/autorun paths)
    crit: notable
    conf: 0.8
    r#for: [javascript, typescript]
    if:
      type: string
      regex: "APPDATA|Startup|crontab|launchd|systemd"

  # File write operations
  - id: file-write-ops
    desc: File write or copy operations
    crit: inert
    conf: 1.0
    r#for: [javascript, typescript]
    if:
      type: string
      regex: "write(File|FileSync)|copy(File|FileSync)"

composite_rules:
  # State helpers
  - id: state-markers
    desc: State management markers (globalState or
    crit: notable
    conf: 0.75
    r#for: [javascript, typescript]
    count_min: 1
    any:
      - id: globalState-word
      - id: postInstall-word

  # Credential paths
  - id: credential-paths
    desc: Sensitive credential file/directory paths
    crit: notable
    conf: 0.8
    r#for: [javascript, typescript]
    count_min: 1
    any:
      - id: ssh-dir
      - id: aws-dir
      - id: npmrc-file
      - id: env-file

  # Send methods
  - id: send-methods
    desc: HTTP send methods
    crit: notable
    conf: 0.7
    r#for: [javascript, typescript]
    count_min: 1
    any:
      - id: post-method
      - id: fetch-call
      - id: axios-lib

  # Post-Install Code Execution
  - id: post-install-exec
    desc: Post-install script execution (dropper pattern)
    crit: suspicious
    conf: 0.95
    mbc: "B0024"
    attack: "T1059"
    r#for: [javascript, typescript]
    all:
      - id: activate-word
      - id: state-markers
      - id: exec-pattern

  # Extension Credential Theft
  - id: credential-theft
    desc: VSCode extension accessing sensitive credentials
    crit: suspicious
    conf: 0.9
    mbc: "E1519"
    attack: "T1555"
    r#for: [javascript, typescript]
    all:
      - id: credential-paths
      - id: readFile-pattern
      - id: send-methods

  # Staged Decryption and Execution inside a VSCode Extension
  - id: staged-dropper
    desc: VSCode extension with staged decryption
    crit: suspicious
    conf: 0.99
    mbc: "B0032"
    attack: "T1027"
    r#for: [javascript, typescript]
    all:
      - id: js-decrypt-eval
      - id: api-usage
      - id: vscode-namespace

  # Extension Persistence via well-known hooks
  - id: persistence
    desc: Extension establishing persistence
    crit: suspicious
    conf: 0.95
    mbc: "F0011"
    attack: "T1547"
    r#for: [javascript, typescript]
    all:
      - id: api-usage
      - id: persistence-locations
      - id: file-write-ops
