# traits/known-tools/backdoor/rustdoor/traits.yaml
# Malware Family - Rustdoor macOS Backdoor
# Rust-based macOS backdoor
# Reference: https://medium.com/s2wblog/rustdoor-and-gatedoor-a-new-pair-of-weapons-disguised-as-legitimate-software-by-suspected-34c94e558b40
# ATT&CK: T1059 (Command and Scripting Interpreter)

traits:
  # === Atomic Traits ===
  # NOTE: botkill is defined in c2/botkill/traits.yaml

  - id: backdoor/rustdoor/ziptask
    description: Rustdoor ziptask dialog
    crit: suspicious
    conf: 0.9
    family: true
    attack: "T1059"
    file_types: [macho]
    condition:
      type: string
      exact: "ziptask"

  - id: backdoor/rustdoor/upload-files
    description: Rustdoor file upload function
    crit: suspicious
    conf: 0.85
    family: true
    attack: "T1041"
    file_types: [macho]
    condition:
      type: string
      exact: "upload_filesrc"

  - id: backdoor/rustdoor/launch-agents
    description: Rustdoor LaunchAgents persistence
    crit: suspicious
    conf: 0.8
    family: true
    attack: "T1543.001"
    file_types: [macho]
    condition:
      type: string
      exact: "LaunchAgents.plist"

  # === V2 config traits ===
  # Note: "daemonize" alone is too generic (matches any daemon)
  # Detection relies on the composite rule requiring multiple config strings

  - id: backdoor/rustdoor/cfg-check-cron
    description: Rustdoor v2 cron check config
    crit: suspicious
    conf: 0.9
    family: true
    attack: "T1053.003"
    file_types: [macho]
    condition:
      type: string
      exact: "check_cron_asked"

  - id: backdoor/rustdoor/cfg-lock-cron
    description: Rustdoor v2 cron lock config
    crit: suspicious
    conf: 0.9
    family: true
    attack: "T1053.003"
    file_types: [macho]
    condition:
      type: string
      exact: "lock_in_cron"

  - id: backdoor/rustdoor/cfg-lock-dock
    description: Rustdoor v2 dock lock config
    crit: suspicious
    conf: 0.9
    family: true
    attack: "T1543.001"
    file_types: [macho]
    condition:
      type: string
      exact: "lock_in_dock"

  - id: backdoor/rustdoor/cfg-lock-launch
    description: Rustdoor v2 launch lock config
    crit: suspicious
    conf: 0.9
    family: true
    attack: "T1543.001"
    file_types: [macho]
    condition:
      type: string
      exact: "lock_in_launch"

  - id: backdoor/rustdoor/cfg-copy-files
    description: Rustdoor v2 file copy config
    crit: suspicious
    conf: 0.5
    family: true
    attack: "T1105"
    file_types: [macho]
    condition:
      type: string
      exact: "copy_files"

composite_rules:
  # Rustdoor v1 detection
  - id: backdoor/rustdoor/v1-detected
    description: Rustdoor macOS backdoor v1
    crit: hostile
    conf: 0.95
    family: true
    reference: "https://medium.com/s2wblog/rustdoor-and-gatedoor-a-new-pair-of-weapons-disguised-as-legitimate-software-by-suspected-34c94e558b40"
    attack: "T1059"
    mbc: "B0022"
    file_types: [macho]
    requires_all:
      - id: c2/botkill/botkill
      - id: backdoor/rustdoor/ziptask
      - id: backdoor/rustdoor/upload-files
      - id: backdoor/rustdoor/launch-agents

  # Rustdoor v2 detection
  - id: backdoor/rustdoor/v2-detected
    description: Rustdoor macOS backdoor v2
    crit: hostile
    conf: 0.95
    family: true
    attack: "T1059"
    mbc: "B0022"
    file_types: [macho]
    requires_count: 4
    any:
      - id: backdoor/rustdoor/cfg-check-cron
      - id: backdoor/rustdoor/cfg-lock-cron
      - id: backdoor/rustdoor/cfg-lock-dock
      - id: backdoor/rustdoor/cfg-lock-launch
      - id: backdoor/rustdoor/cfg-copy-files
