# traits/lateral/supply-chain/vscode/malicious.yaml
# Malicious VSCode Extension Patterns
# Detects high-confidence malicious patterns in VSCode extensions
# These combine multiple indicators to reduce false positives

traits:
  # VSCode-specific markers for contextual analysis
  - id: api-usage
    description: Uses VSCode API (extension marker)
    criticality: inert
    confidence: 1.0
    file_types: [javascript, typescript]
    condition:
      type: yara
      source: |
        rule vscode_api_usage {
          strings:
            $require = /require\(['"]vscode['"]\)/ ascii
            $import = /from\s+['"]vscode['"]/ ascii
            $namespace = "vscode." ascii
            $activate = "activate" ascii fullword
          condition:
            any of them
        }

  # Post-Install Code Execution
  - id: post-install-exec
    description: Post-install script execution (dropper pattern)
    criticality: suspicious
    confidence: 0.95
    mbc: "B0024"
    attack: "T1059"
    file_types: [javascript, typescript]
    condition:
      type: yara
      source: |
        rule vscode_post_install_exec {
          strings:
            $activate = "activate" ascii fullword
            $state1 = "globalState" ascii
            $state2 = "postInstall" ascii
            $exec = /exec(Sync|Async)?/ ascii
          condition:
            $activate and any of ($state*) and $exec
        }

  # Extension Credential Theft
  - id: credential-theft
    description: VSCode extension accessing sensitive credentials
    criticality: suspicious
    confidence: 0.9
    mbc: "E1519"
    attack: "T1555"
    file_types: [javascript, typescript]
    condition:
      type: yara
      source: |
        rule vscode_credential_theft {
          strings:
            $cred1 = ".ssh" ascii
            $cred2 = ".aws" ascii
            $cred3 = ".npmrc" ascii
            $cred4 = ".env" ascii
            $read = /read(File|FileSync)/ ascii
            $send1 = "POST" ascii
            $send2 = "fetch(" ascii
            $send3 = "axios" ascii
          condition:
            any of ($cred*) and $read and any of ($send*)
        }

  # VSCode namespace usage pattern
  - id: vscode-namespace
    description: Uses VSCode namespace
    criticality: inert
    confidence: 1.0
    file_types: [javascript, typescript]
    condition:
      type: string
      regex: "vscode\\."

  # Persistence location references
  - id: persistence-locations
    description: References persistence locations (startup/autorun paths)
    criticality: notable
    confidence: 0.8
    file_types: [javascript, typescript]
    condition:
      type: string
      regex: 'APPDATA|Startup|crontab|launchd|systemd'

  # File write operations
  - id: file-write-ops
    description: File write or copy operations
    criticality: inert
    confidence: 1.0
    file_types: [javascript, typescript]
    condition:
      type: string
      regex: 'write(File|FileSync)|copy(File|FileSync)'

composite_rules:
  # Staged Decryption and Execution inside a VSCode Extension
  - id: staged-dropper
    description: VSCode extension with staged decryption and execution (Dropper)
    criticality: hostile
    confidence: 0.99
    mbc: "B0032"
    attack: "T1027"
    file_types: [javascript, typescript]
    requires_all:
      - id: js-decrypt-eval
      - id: api-usage
      - id: vscode-namespace

  # Extension Persistence via well-known hooks
  - id: persistence
    description: Extension establishing persistence
    criticality: hostile
    confidence: 0.95
    mbc: "F0011"
    attack: "T1547"
    file_types: [javascript, typescript]
    requires_all:
      - id: api-usage
      - id: persistence-locations
      - id: file-write-ops
